# curator.cron
1 * * * * LC_ALL=C.UTF-8 /usr/bin/curator --config /tmp/config.yaml /tmp/action.yaml &>> /var/log/cron.log

